# dike.service
[Unit]
Description=DIKE - Stake.com Dice Auto Bet Bot
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/dike
ExecStart=/usr/bin/python3 /root/dike/bot.py
Restart=always
RestartSec=10
StandardOutput=append:/root/dike/dike.log
StandardError=append:/root/dike/dike.log

[Install]
WantedBy=multi-user.target
